# files/system/usr/lib/systemd/system/cipher-cleaner.service
[Unit]
Description=Cipherblue Cleaning Service
DefaultDependencies=no
Before=sysinit.target
Requires=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ProtectSystem=full
NoNewPrivileges=yes
TimeoutStartSec=0
RemainAfterExit=yes

# 1. Safely clean temporary shared memory and mqueue
ExecStart=/usr/bin/find /dev/shm/ -mindepth 1 -delete
ExecStart=/usr/bin/find /dev/mqueue/ -mindepth 1 -delete
ExecStart=/usr/bin/find /var/tmp/ -mindepth 1 -delete

# ==============================================================================
# CIPHERBLUE GITOPS HOSTS INJECTOR
# ==============================================================================
# Extracts the static blocklist baked by the GitHub Runner and cleanly injects
# it into the live /etc/hosts before the network stack is initialized.
ExecStart=/usr/bin/sh -c 'sed -i "/# CIPHERBLUE PRIVATE DOMAIN BLOCKLIST/,$d" /etc/hosts && echo "# CIPHERBLUE PRIVATE DOMAIN BLOCKLIST" >> /etc/hosts && cat /etc/cipherblue/hosts.blocklist >> /etc/hosts || true'

# ==============================================================================
# CIPHERBLUE GARBAGE COLLECTOR (DETERMINISTIC ENGINE)
# ==============================================================================
ExecStart=-/usr/bin/sh -c 'while IFS=: read -r user pass uid gid info home shell; do \
    if [ "$uid" -eq 0 ] || [ "$uid" -ge 1000 ]; then \
        if [ -d "$home/.local/share/flatpak" ]; then \
            rm -rf "$home/.local/share/flatpak"; \
        fi; \
    fi; \
done < /etc/passwd'

StandardOutput=null
StandardError=null

[Install]
WantedBy=sysinit.target