# files/system/usr/lib/systemd/system/cipher-secure-rebase.service
[Unit]
Description=Cipherblue Cryptographic Signature Enforcer
After=network-online.target rpm-ostreed.service
Wants=network-online.target

# CRITICAL: This physically stops the service from running on every boot 
# once the machine has been successfully rebased to the signed image.
ConditionPathExists=!/var/lib/cipherblue-signed-rebase.stamp

[Service]
Type=simple
ExecStart=/usr/bin/bash /usr/libexec/cipherblue/cipher-secure-rebase.sh

# If the network drops or rpm-ostree locks up, retry gracefully every 60 seconds
Restart=on-failure
RestartSec=60s
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target