# files/system/usr/lib/systemd/system/cipher-cleanup.service
[Unit]
Description=Cipherblue Deep Cleanup Engine
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/cipherblue/cipher-cleanup.sh