# files/system/usr/lib/systemd/system/cipher-audit-drift.service
[Unit]
Description=Cipherblue Configuration Drift Auditor
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/cipherblue/cipher-audit-drift.sh