# files/system/usr/lib/systemd/system/cipher-user-env-lockdown.service
[Unit]
Description=Cipherblue User Environment Immutability Engine
# CRITICAL: Prevents Systemd from automatically injecting After=sysinit.target
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/cipherblue/cipher-user-env-lockdown.sh
RemainAfterExit=yes

[Install]
# We bind it directly to sysinit instead of multi-user so it runs at the kernel level
WantedBy=sysinit.target