# files/system/usr/lib/systemd/system/cipher-flatpak-update.service
[Unit]
Description=Cipherblue Automated App Updater
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/cipherblue/cipher-flatpak-update.sh