# files/system/usr/lib/systemd/system/cipher-cleaner.service
[Unit]
Description=Cipherblue Cleaning Service
DefaultDependencies=no
Before=sysinit.target
Requires=local-fs.target
After=local-fs.target

[Service]
Type=oneshot
ProtectSystem=full

# HARDENING: Punch a micro-hole in ProtectSystem=full just for this exact file.
ReadWritePaths=/etc/hosts

NoNewPrivileges=yes
TimeoutStartSec=0
RemainAfterExit=yes

# 1. Safely clean temporary shared memory and mqueue
ExecStart=/usr/bin/find /dev/shm/ -mindepth 1 -delete
ExecStart=/usr/bin/find /dev/mqueue/ -mindepth 1 -delete
ExecStart=/usr/bin/find /var/tmp/ -mindepth 1 -delete

# ==============================================================================
# CIPHERBLUE GITOPS HOSTS INJECTOR (Inode-Preserving Engine)
# ==============================================================================
ExecStart=/usr/bin/sh -c "sed '/# CIPHERBLUE PRIVATE DOMAIN BLOCKLIST/,$$d' /etc/hosts > /run/hosts.clean && echo '' >> /run/hosts.clean && echo '# CIPHERBLUE PRIVATE DOMAIN BLOCKLIST' >> /run/hosts.clean && cat /etc/cipherblue/hosts.blocklist >> /run/hosts.clean && cat /run/hosts.clean > /etc/hosts && rm -f /run/hosts.clean"

# ==============================================================================
# CIPHERBLUE GARBAGE COLLECTOR (Surgical Engine)
# ==============================================================================
# CRITICAL FIX 5: We stop "carpet bombing" the flatpak directory, which destroyed 
# Flatseal's 'overrides' vault. We now surgically annihilate unauthorized apps 
# and repos, while preserving and self-healing the overrides directory!
ExecStart=-/usr/bin/sh -c 'while IFS=: read -r user pass uid gid info home shell; do \
    if [ "$$uid" -eq 0 ] || [ "$$uid" -ge 1000 ]; then \
        if [ -d "$$home/.local/share/flatpak" ]; then \
            rm -rf "$$home/.local/share/flatpak/app" "$$home/.local/share/flatpak/runtime" "$$home/.local/share/flatpak/repo" "$$home/.local/share/flatpak/exports"; \
        fi; \
        mkdir -p "$$home/.local/share/flatpak/overrides"; \
        chown -R "$$uid:$$gid" "$$home/.local/share/flatpak" 2>/dev/null || true; \
    fi; \
done < /etc/passwd'

StandardOutput=null
StandardError=null

[Install]
WantedBy=sysinit.target
