# files/system/etc/systemd/system/cipher-remount.service.d/cipherblue.conf
[Service]
NoNewPrivileges=yes
PrivateTmp=yes
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077

