# files/system/etc/systemd/system/systemd-localed.service.d/cipherblue.conf
[Service]
NoNewPrivileges=yes
PrivateTmp=yes
ProtectClock=true
ProtectControlGroups=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077

