# files/system/usr/lib/systemd/user/cipherblue-firmware-check.service
[Unit]
Description=cipherblue Firmware Update Notifications
Wants=graphical-session.target
After=graphical-session.target
ConditionUser=!@system
ConditionPathIsDirectory=/var/home/%u

[Service]
Type=simple
ExecCondition=/usr/libexec/cipherblue/firmwarecheckoutofdate
ExecStart=/usr/libexec/cipherblue/firmwareoutofdatenotify
