# files/system/usr/lib/systemd/user/cipherblue-firmware-check.timer
[Unit]
Description=cipherblue Firmware Update Checks
ConditionUser=!@system

[Timer]
RandomizedDelaySec=1m
OnStartupSec=2m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
