# files/system/usr/lib/systemd/user/security-update-notification.path
[Unit]
Description=Monitor for rpm-ostreed-automatic completion
Wants=paths.target
ConditionUser=!@system

[Path]
PathChanged=/var/lib/secureblue/rpm-ostreed-automatic-completed.stamp

[Install]
WantedBy=paths.target
