# files/system/usr/lib/systemd/user/secureblue-update-verification.timer
[Unit]
Description=secureblue update verification timer
ConditionUser=!@system

[Timer]
RandomizedDelaySec=1m
OnStartupSec=2m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
