# files/system/usr/lib/systemd/user/secureblue-flatpak-setup.timer
[Unit]
Description=secureblue flatpak setup

[Timer]
RandomizedDelaySec=1s
OnBootSec=2s
Persistent=true

[Install]
WantedBy=timers.target
