# files/system/etc/systemd/user/org.gnome.Shell@user.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gnome-shell --mode=%i --no-x11
