# files/system/etc/systemd/system/flatpak-system-update.service.d/99-cipherblue-cleanup.conf
[Service]
# The '-' prefix tells systemd not to fail the service if no unused flatpaks are found.
# 1. Clean up unused system runtimes left behind by upstream updates
ExecStartPost=-/usr/bin/flatpak --system uninstall --unused -y --noninteractive --delete-data

# 2. Obliterate unauthorized user flatpaks, repos, and their data
ExecStartPost=-/usr/bin/sh -c 'rm -rf /var/home/*/.local/share/flatpak'