# files/system/etc/systemd/system/user@.service.d/99-cipher-run-hardening.conf
[Service]
# The '+' tells systemd to run this specific command as root, even though the service belongs to the user.
# The '-' tells systemd to ignore it if it fails (e.g., if the user has no runtime dir).
# %i dynamically injects the UID (e.g., 1000) of the user logging in.
ExecStartPre=+-/usr/bin/mount -o remount,nosuid,noexec,nodev /run/user/%i