# files/system/usr/lib/systemd/system/user-runtime-dir@.service.d/10-noexec.conf
[Unit]
After=local-fs.target

[Service]
# Ensure per-user runtime dirs get noexec,nosuid,nodev (ignore errors for safety)
ExecStartPost=-/usr/libexec/cipherblue/remount-run-user %i
