# files/system/etc/profile.d/hardened_malloc.sh
#!/usr/bin/sh
export LD_PRELOAD='libhardened_malloc.so libno_rlimit_as.so'