# files/system/usr/lib/systemd/system/user@.service.d/10-restrict-delegation.conf
[Service]
# Disallow cgroup delegation to user instances to block rootless containers
Delegate=no
# Reasonable per-user cap on number of tasks
TasksMax=4096

