# files/system/etc/systemd/system/fwupd-refresh.service.d/cipherblue.conf
[Service]
NoNewPrivileges=yes
PrivateTmp=yes
ProtectClock=true
ProtectControlGroups=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
UMask=0077

