# files/system/usr/lib/systemd/system/dnsconfd.service.d/secureblue.conf
# Copyright 2025 The Secureblue Authors
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

[Service]
# SELinux policy allows net_bind_service net_raw net_admin, supposedly for
# unbound-control.
CapabilityBoundingSet=
ProtectSystem=strict
ReadWritePaths=/etc/resolv.conf /run/dnsconfd
ProtectHome=yes
PrivateTmp=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectProc=invisible
LockPersonality=yes
RestrictNamespaces=yes
PrivateDevices=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
PrivateMounts=yes
IPAddressDeny=any
IPAddressAllow=localhost
RestrictRealtime=yes
