# files/system/kinoite/etc/systemd/user/plasma-xembedsniproxy.service.d/override.conf
[Service]
# Do not start service if Xwayland is disabled
ExecCondition=/usr/bin/bash -c 'systemctl --user show --property=ExecStart plasma-kwin_wayland.service | grep -q -e "--xwayland"'
