# files/system/desktop/usr/lib/systemd/user/security-update-notification.service
[Unit]
Description=Security update notification
Wants=graphical-session.target
After=graphical-session.target
ConditionUser=!@system

[Service]
Type=simple
ExecStart=/usr/libexec/secureblue/security-update-notification
