# files/system/desktop/usr/lib/systemd/user/secureblue-flatpak-setup.timer
[Unit]
Description=secureblue flatpak setup

[Timer]
OnStartupSec=5s

[Install]
WantedBy=timers.target
