# files/system/etc/systemd/system/systemd-resolved.service.d/cipherblue.conf
[Service]
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelTunables=true
ProtectProc=invisible
MemoryDenyWriteExecute=true
RestrictNamespaces=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
UMask=0077
LockPersonality=true
CapabilityBoundingSet=
SystemCallFilter=@system-service ~@privileged ~@module ~@mount ~@keyring ~@resources ~@raw-io
