# files/system/etc/systemd/user/xdg-user-dirs.service.d/override.conf
[Service]
# Clear the upstream execution command entirely
ExecStart=
# Replace with a mathematical no-op. 
# Because our Immutability Engine pre-creates all XDG directories flawlessly,
# this background updater is obsolete and causes atomic-save crashes on our frozen nodes.
ExecStart=/usr/bin/true