# files/system/etc/systemd/system/rpm-ostreed-automatic.timer.d/99-cipherblue-aggressive.conf
[Unit]
# Ensure the timer clock does not start until the network is completely established
Wants=network-online.target
After=network-online.target

[Timer]
# Step 1: Defensive clear of any lingering upstream schedules
OnBootSec=
OnCalendar=
OnUnitActiveSec=

# Step 2: Define the Network-Aware Aggressive Schedule
# Trigger exactly 1 minute AFTER the network connects
OnActiveSec=1min

# Check for zero-day patches every 4 hours while the computer is in use
OnUnitActiveSec=4h

# A tiny delay prevents server throttling if you deploy multiple appliances
RandomizedDelaySec=1min

# If the system missed a schedule, run immediately
Persistent=true