# files/system/etc/xdg/autostart/cipherblue-mute-mic.desktop
[Desktop Entry]
Type=Application
Name=Script to Mute Microphone on Boot
Exec=/usr/bin/amixer set Capture nocap
Icon=utilities-terminal
